# services/agriha-ui.service
[Unit]
Description=AgriHA Web UI (Dashboard + Settings)
After=network.target agriha-daemon.service
Wants=agriha-daemon.service

[Service]
Type=simple
User=unipi
WorkingDirectory=/opt/agriha-ui
ExecStart=/opt/agriha-ui/venv/bin/uvicorn agriha_ui.app:app --host 0.0.0.0 --port 8502
Restart=on-failure
RestartSec=5
Environment=UNIPI_API_URL=http://localhost:8080
Environment=SYSTEM_PROMPT_PATH=/etc/agriha/system_prompt.txt
Environment=AGRIHA_THRESHOLDS_PATH=/etc/agriha/thresholds.yaml
Environment=CONTROL_LOG_DB=/var/lib/agriha/control_log.db

[Install]
WantedBy=multi-user.target
