# systemd/agriha-linebot.service
# /etc/systemd/system/agriha-linebot.service
# LINE Bot (FastAPI/uvicorn) — VPS systemdサービス
#
# デプロイ:
#   sudo cp systemd/agriha-linebot.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now agriha-linebot

[Unit]
Description=AgriHA LINE Bot (FastAPI + uvicorn)
After=network.target wg-quick@wg-farmers.service
Wants=wg-quick@wg-farmers.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/agriha-linebot
EnvironmentFile=/opt/agriha-linebot/.env
Environment="PATH=/opt/agriha-linebot/.venv/bin:/usr/local/bin:/usr/bin"
ExecStart=/opt/agriha-linebot/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 8443
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agriha-linebot

[Install]
WantedBy=multi-user.target
