# systemd/agriha-ui.service
[Unit]
Description=AgriHA Local Web UI (FastAPI + htmx)
After=network.target unipi-daemon.service

[Service]
Type=simple
User=agriha
Group=agriha
WorkingDirectory=__REPO_DIR__
EnvironmentFile=__REPO_DIR__/.env
Environment="PATH=__REPO_DIR__/.venv/bin:/usr/local/bin:/usr/bin"
ExecStart=__REPO_DIR__/.venv/bin/uvicorn agriha.chat.app:app --host 0.0.0.0 --port 8501
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agriha-ui

[Install]
WantedBy=multi-user.target
