# systemd/agriha-llm.service
[Unit]
Description=AgriHA LLM Server (llama-server + LFM2.5)
After=network.target
Documentation=https://github.com/ggml-org/llama.cpp

[Service]
Type=simple
User=root
Environment=LD_LIBRARY_PATH=/opt/llama-server/bin
ExecStart=/opt/llama-server/bin/llama-server \
    -m /opt/llama-server/models/LFM2.5-1.2B-Instruct-Q4_K_M.gguf \
    --port 8081 \
    -c 4096 \
    -t 4 \
    --mlock \
    --jinja
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agriha-llm

[Install]
WantedBy=multi-user.target
