# systemd/agriha-nullclaw-proxy.service
[Unit]
Description=NullClaw OpenAI-Compatible Proxy
After=network.target

[Service]
Type=simple
User=agriha
ExecStart=/usr/local/bin/uvicorn agriha.control.nullclaw_proxy:app --host 127.0.0.1 --port 3001
WorkingDirectory=/opt/agriha
Restart=always
RestartSec=5
EnvironmentFile=/opt/agriha/.env

[Install]
WantedBy=multi-user.target
