# systemd/agriha-ui.service
[Unit]
Description=AgriHA Local Web UI (FastAPI + htmx)
After=network.target unipi-daemon.service

[Service]
Type=simple
User=agriha
Group=agriha
WorkingDirectory=/home/agriha/uecs-llm
EnvironmentFile=/home/agriha/uecs-llm/.env
Environment="PATH=/home/agriha/uecs-llm/venv/bin:/usr/local/bin:/usr/bin"
ExecStart=/home/agriha/uecs-llm/venv/bin/uvicorn agriha_ui.app:app --host 0.0.0.0 --port 8502
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agriha-ui

[Install]
WantedBy=multi-user.target
